# webapi/deploy/com.opentune.web.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job that keeps the OpenTune web API running and restarts it after a crash or reboot.

  Setup:
    1. Edit the REPLACE_ME paths below (repo location, your username, a JDK 21 home).
    2. Copy to ~/Library/LaunchAgents/com.opentune.web.plist
    3. launchctl load -w ~/Library/LaunchAgents/com.opentune.web.plist
  Stop:   launchctl unload -w ~/Library/LaunchAgents/com.opentune.web.plist
  Logs:   tail -f ~/Library/Logs/opentune-web.log
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.opentune.web</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>REPLACE_ME/OpenTune/webapi/deploy/run.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- A Java 21+ runtime. Homebrew: /opt/homebrew/opt/openjdk@21. Or the Android Studio JBR. -->
        <key>OPENTUNE_JAVA_HOME</key>
        <string>REPLACE_ME_JDK21_HOME</string>
        <key>PORT</key>
        <string>8080</string>
    </dict>

    <!-- Keep it running and bring it back on login/reboot. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>REPLACE_ME_HOME/Library/Logs/opentune-web.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_ME_HOME/Library/Logs/opentune-web.log</string>
</dict>
</plist>
