# webapi/deploy/com.opentune.web.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job that keeps the OpenTune web API running and restarts it after a crash or reboot.

  It runs the STAGED launcher (~/opentune/server/...), not the repo, because launchd agents cannot
  read files under ~/Desktop, ~/Documents or ~/Downloads on macOS. Run webapi/deploy/stage.sh to
  create that copy. If your repo is not under a protected folder, you may point ProgramArguments at
  REPO/webapi/build/install/webapi/bin/webapi directly and drop OPENTUNE_WEB_ROOT.

  Setup:
    1. build.sh, then stage.sh (creates ~/opentune).
    2. Edit the REPLACE_ME values below.
    3. Copy to ~/Library/LaunchAgents/com.opentune.web.plist
    4. launchctl load -w ~/Library/LaunchAgents/com.opentune.web.plist
  Stop:   launchctl unload -w ~/Library/LaunchAgents/com.opentune.web.plist
  Logs:   tail -f ~/Library/Logs/opentune-web.log
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.opentune.web</string>

    <key>ProgramArguments</key>
    <array>
        <string>REPLACE_ME_HOME/opentune/server/bin/webapi</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- A Java 21+ runtime. Homebrew: /opt/homebrew/opt/openjdk@21. Or the Android Studio JBR. -->
        <key>JAVA_HOME</key>
        <string>REPLACE_ME_JDK21_HOME</string>
        <key>PORT</key>
        <string>8080</string>
        <!-- The staged web bundle, so the server does not reach back into the repo. -->
        <key>OPENTUNE_WEB_ROOT</key>
        <string>REPLACE_ME_HOME/opentune/web-dist</string>
        <!-- Absolute path to yt-dlp. REQUIRED under launchd: its PATH is minimal (no /opt/homebrew/bin),
             so a bare "yt-dlp" is not found and the stream-resolution fallback silently never runs --
             which leaves any track InnerTube cannot resolve unplayable. Find yours with `which yt-dlp`. -->
        <key>OPENTUNE_YTDLP_PATH</key>
        <string>REPLACE_ME_YTDLP_PATH</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>REPLACE_ME_HOME/opentune</string>

    <!-- Keep it running and bring it back on login/reboot. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>REPLACE_ME_HOME/Library/Logs/opentune-web.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_ME_HOME/Library/Logs/opentune-web.log</string>
</dict>
</plist>
